# pgpool/4.5.2 Ubuntu 22.04/fix socket glitch on Ubuntu, add vip automatic mount/pgpool_repair service & timer/pgpool_repair.timer
# Timer for the service

[Unit]
Description=Executes socket repair operation
Requires=pgpool_repair.service

[Timer]
Unit=pgpool_repair.service
OnCalendar=*:*:0/10
AccuracySec=1s

[Install]
WantedBy=timers.target

