# pgpool/4.5.2 Ubuntu 22.04/fix socket glitch on Ubuntu/pgpool_repair service & timer/pgpool_repair.timer
# Timer for the service

[Unit]
Description=Executes socket repair operation
Requires=pgpool_socket_repair@9999.service

[Timer]
Unit=pgpool_socket_repair@9999.service
OnCalendar=*:*:0/10
AccuracySec=1s

[Install]
WantedBy=timers.target

